# guardian/crontab
0 */12 * * * /usr/local/bin/python3 ../app/main.py >> cron_log.log 2>&1 
